# dep/etc/profile.d/S80RkEnv.sh
export HOME=/oem
export PATH=$PATH:$HOME:$HOME/bin:$HOME/usr/bin:$HOME/sbin:$HOME/usr/sbin:$HOME/usr/local/bin:/usr/local/bin
export LD_LIBRARY_PATH=$HOME/usr/lib:$HOME/lib:$LD_LIBRARY_PATH

alias ll="ls -al"

